# docker/webgrinder/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7.3]
command=/bin/sh -c "mkdir /run/php && /usr/sbin/php-fpm7.3 --nodaemonize"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"